# build/config/mainrouter/files/etc/uci-defaults/9961-network-custom
set network.globals.packet_steering='1'

/etc/init.d/network restart
/etc/init.d/odhcpd reload
