# build/config/rpi4b/files/etc/uci-defaults/9901-security
passwd root << EOF
${ROOT_PASSWORD_TEMPLATE}
${ROOT_PASSWORD_TEMPLATE}
EOF

# TODO: Let's also provide dropbear_ed25519_host_key  dropbear_rsa_host_key
# so they are kept between the builds?
cat >> /etc/dropbear/authorized_keys << EOI
${SSH_AUTHORIZED_KEYS_TEMPLATE}
EOI

uci set dropbear.@dropbear[0].Interface='lan'
uci set dropbear.@dropbear[0].PasswordAuth='off'
uci set dropbear.@dropbear[0].RootPasswordAuth='off'
uci commit dropbear

/etc/init.d/dropbear restart
