# build/config/rpi4b/files/etc/uci-defaults/9901-security
passwd root << EOF
{{ .Env.ROOT_PASSWORD_TEMPLATE }}
{{ .Env.ROOT_PASSWORD_TEMPLATE }}
EOF

cat >> /etc/dropbear/authorized_keys << EOI
{{ .Env.SSH_AUTHORIZED_KEYS_TEMPLATE }}
EOI

uci batch << EOI
set dropbear.@dropbear[0].Interface='lan'
set dropbear.@dropbear[0].PasswordAuth='off'
set dropbear.@dropbear[0].RootPasswordAuth='off'
commit dropbear
EOI

/etc/init.d/dropbear restart
