# build/config/rpi4b/files/etc/uci-defaults/9905-dhcp
uci batch << EOI
# This is going to be changed by https-dns-proxy in the end
# Cloudflare is just a backup
set dhcp.@dnsmasq[0].server='1.1.1.1'

set dhcp.wan.start='100'
set dhcp.wan.limit='150'
set dhcp.wan.leasetime='12h'
set dhcp.wan.ra_flags='none'

set dhcp.guest=dhcp
set dhcp.guest.interface='guest'
set dhcp.guest.start='100'
set dhcp.guest.limit='150'
set dhcp.guest.leasetime='12h'

${DHCP_CUSTOM_HOSTS_CONFIG_TEMPLATE}

commit dhcp
EOI

/etc/init.d/dnsmasq reload
/etc/init.d/odhcpd reload
