# build/config/rpi4b/files/etc/uci-defaults/9908-wireguard
WG_IF="${WIREGUARD_INTERFACE_TEMPLATE}"

uci batch << EOI
set network.${WG_IF}=interface
set network.${WG_IF}.proto='wireguard'
set network.${WG_IF}.force_link='1'
set network.${WG_IF}.private_key='${VPN_PRIVATE_KEY_TEMPLATE}'
set network.${WG_IF}.listen_port='${VPN_LISTEN_PORT_TEMPLATE}'
set network.${WG_IF}.addresses='${VPN_ADDRESSES_TEMPLATE}'
commit network
EOI

${VPN_CUSTOM_CLIENTS_CONFIG_TEMPLATE}

uci commit network

/etc/init.d/network restart
