# build/config/rpi4b/files/etc/uci-defaults/9912-mdns-repeater
set -e

# TODO: we need repeating from IoT network when it's added
# probably the eth2 for guest network can be removed then?

uci del mdns_repeater.main.interface
uci add_list mdns_repeater.main.interface='eth2'
uci add_list mdns_repeater.main.interface='br-lan'
uci commit

/etc/init.d/mdns-repeater restart
