# build/config/rpi4b/files/etc/uci-defaults/9912-mdns-repeater
set -e

# We only want the IoT and LAN to be able to see each other
# So not listing the guest VLAN here
uci del mdns_repeater.main.interface || true
uci add_list mdns_repeater.main.interface='br-lan.30'
uci add_list mdns_repeater.main.interface='br-lan.99'
uci commit

/etc/init.d/mdns-repeater restart
