# build/config/rpi4b/files/etc/uci-defaults/9951-irqbalance
uci set irqbalance.irqbalance.enabled='1'
uci commit irqbalance

/etc/init.d/irqbalance start
