# build/config/rpi4b/files/etc/uci-defaults/9953-adblock
uci batch << EOI
set adblock.global.adb_maxqueue='4'
set adblock.global.adb_report='1'
set adblock.global.adb_repiface='br-lan'

del adblock.global.adb_sources
{{- range .adblock.sources }}
add_list adblock.global.adb_sources="{{ . }}"
{{- end }}

commit adblock
EOI

# Full restart needed because of reporting interface
/etc/init.d/adblock restart
