# build/config/rpi4b/files/etc/uci-defaults/9955-mdns-repeater
set -e

uci del mdns_repeater.main.interface || true
{{- range $name, $config := .networks }}
{{- if has $config "broadcast_mdns" }}
{{- if $config.broadcast_mdns }}
{{- if has $config "vlan" }}
uci add_list mdns_repeater.main.interface='{{ $config.device }}.{{ $config.vlan.id }}'
{{ else }}
uci add_list mdns_repeater.main.interface='{{ $config.device }}'
{{- end }}
{{- end }}
{{- end }}
{{- end }}

uci commit mdns_repeater

/etc/init.d/mdns-repeater restart
