# build/config/rpi4b/files/etc/uci-defaults/9962-fan
# Poe hat fan speeds tuning
cat << EOI >> /boot/config.txt
dtparam=poe_fan_temp0=50000
dtparam=poe_fan_temp1=60000
dtparam=poe_fan_temp2=70000
dtparam=poe_fan_temp3=80000
EOI

# TODO: we need to reboot after this operation
# can we do it from here?
#
# TODO: /boot seems to persist across `firstboot` commands
# we may need a bit more sophisticated approach than just appending to the end of the file
# Probably the best option would be to bake this in the image directly during the build
# Maybe by patching the `target/linux/bcm27xx/image/config.txt` in openwrt build dir,
# before actually starting the build?
