# build/rpi4b/files/etc/uci-defaults/9907-adblock
uci -q batch << EOI
set adblock.global.adb_maxqueue='4'
set adblock.global.adb_report='1'
set adblock.global.adb_repiface='br-lan'

${ADBLOCK_SOURCES_TEMPLATE}

commit adblock
EOI

# Full restart needed because of reporting interface
/etc/init.d/adblock restart
