# build/config/asus-rt-ax53u/files/etc/uci-defaults/9901-security
passwd root << EOF
{{ .Env.ROOT_PASSWORD_TEMPLATE }}
{{ .Env.ROOT_PASSWORD_TEMPLATE }}
EOF

cat >> /etc/dropbear/authorized_keys << EOI
{{ .Env.SSH_AUTHORIZED_KEYS_TEMPLATE }}
EOI

uci -q batch << EOI
set dropbear.@dropbear[0].PasswordAuth='off'
set dropbear.@dropbear[0].RootPasswordAuth='off'
commit dropbear
EOI
