# build/asus-rt-ax53u/files/etc/uci-defaults/9900-base-custom-config
set -e

if [ "${ENABLE_LETSENCRYPT_CERTS_TEMPLATE}" = "true" ]; then
    uci set system.@system[0].hostname='${HOSTNAME_TEMPLATE}.${LETSENCRYPT_HOSTNAME_TEMPLATE}'
    uci commit system
else
    uci set system.@system[0].hostname='${HOSTNAME_TEMPLATE}'
    uci commit system
fi

uci batch << EOI
set system.@system[0].zonename='${ZONENAME_TEMPLATE}'
set system.@system[0].timezone='${TIMEZONE_TEMPLATE}'
set system.@system[0].log_proto='udp'
set system.@system[0].conloglevel='8'
set system.@system[0].cronloglevel='5'
commit system

set uhttpd.main.redirect_https='on'
commit uhttpd
EOI
